# resources/systemd/http2k.service
[Unit]
Description=http2k daemon

[Service]
User=root
ExecStart=/usr/bin/n2kafka /etc/http2k/config.json

[Install]
WantedBy=multi-user.target
